# config/includes.chroot/etc/profile.d/xdg_cache_home.sh
#!/bin/bash
#elaborado por sinfallas
export XDG_CACHE_HOME="/tmp/.cache
